# site-client/crontab
*/30 * * * * /usr/local/bin/docker-compose -f /opt/miireport/docker-compose-yml up > /dev/null
